# debian/prerm
#!/bin/bash
set -e

# stop django server
invoke-rc.d django stop

